# supervisor.conf
[program:{{queuename}}]
command=php {{artisanroot}} queue:work
autostart=true
autorestart=true
stderr_logfile=/var/log/{{queuename}}.err.log
stdout_logfile=/var/log/{{queuename}}.out.log